# debian/inzown-btn.prerm

systemctl stop inzown-btn
systemctl disable inzown-btn

#DEBHELPER#